# files/supervisord/nmbd.conf
[program:nmbd]
command=/usr/local/samba/sbin/nmbd -i
autostart=true
autorestart=true
startsecs=5

