# files/supervisord/smbd.conf
[program:smbd]
command=/usr/local/samba/sbin/smbd -i
autostart=true
autorestart=true
startsecs=5

